# templates/default/couchpotato.service.erb
[Unit]
Description=CouchPotato application instance
After=network.target

[Service]
ExecStart=/home/apps/couchpotato/CouchPotato.py --daemon --data_dir "/home/apps/.config/couchpotato"
GuessMainPID=no
Type=forking
User=apps
Group=apps

[Install]
WantedBy=multi-user.target
