# templates/default/nzbget.service.erb
[Unit]
Description=NZBGet
# This will load nzbget only after setting up the network
After=network.target
[Service]
User=apps
Group=apps
Type=forking
ExecStart=/home/apps/nzbget/nzbget -c /home/apps/.config/nzbget/nzbget.conf -D
ExecStop=/home/apps/nzbget/nzbget -Q
ExecReload=/home/apps/nzbget/nzbget -O
KillMode=process
Restart=on-failure
[Install]
# Needed to run nzbget at boot
WantedBy=multi-user.target
