# templates/default/qbittorrent.service.erb
[Unit]
Description=qbittorrent torrent server

[Service]
User=apps
ExecStart=/usr/bin/qbittorrent-nox
Restart=on-abort

[Install]
WantedBy=multi-user.target
